# systemd/telemux.service
[Unit]
Description=Private Dev Factory control plane
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/private-dev-factory
EnvironmentFile=%h/private-dev-factory/.env
ExecStart=%h/.bun/bin/bun run src/main.ts
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
